# kalki-os-structured/src/core/iso-profile/kalki-base/airootfs/etc/skel/.config/autostart/wine-setup.desktop
[Desktop Entry]
Type=Application
Name=Wine Setup
Exec=/usr/bin/env bash /usr/share/kalki/scripts/wine-setup.sh
X-GNOME-Autostart-enabled=true
NoDisplay=true 