# kalki-os-structured/src/core/iso-profile/kalki-base/airootfs/etc/systemd/user/kalki-apps.service
[Unit]
Description=Kalki OS Application Management Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/opt/kalki/scripts/app-startup-manager.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
