# kalki-os-structured/src/core/iso-profile/kalki-base/airootfs/etc/systemd/system/kalki-ai-setup.service
[Unit]
Description=Kalki OS AI Model Configuration
After=network-online.target
Wants=network-online.target
ConditionFirstBoot=yes

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kalki-ai
RemainAfterExit=yes
StandardOutput=journal

[Install]
WantedBy=multi-user.target
