# kalki-os-structured/src/core/iso-profile/kalki-base/airootfs/etc/systemd/system/kalki-install-assistant.service
[Unit]
Description=Kalki Installation Assistant
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/kalki/ai-assistant/kia.py --install
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
