# kalki-os-structured/src/core/iso-profile/kalki-base/airootfs/etc/systemd/system/omnet.service
[Unit]
Description=OMNet Neural Core - Kalki OS AI Coordination
After=network.target ollama.service
Requires=ollama.service

[Service]
Type=simple
User=kalki
Group=kalki
WorkingDirectory=/opt/kalki/omnet
ExecStart=/usr/bin/python3 /opt/kalki/omnet/omnet_core.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/kalki

[Install]
WantedBy=multi-user.target
