# kalki-os-restructured/src/core/iso-profile/kalki-base/airootfs/etc/systemd/system/kalki-wallpaper.timer
[Unit]
Description=Update Kalki OS Wallpaper Based on Time of Day

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
