# kalki-os-restructured/src/core/iso-profile/kalki-base/airootfs/etc/systemd/system/kalki-gesture.service
[Unit]
Description=Kalki OS Gesture Recognition Service
After=network.target graphical.target
Requires=omnet.service

[Service]
Type=simple
User=kalki
Group=video
Environment=DISPLAY=:0
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/kalki/gesture
ExecStart=/usr/local/bin/kalki-gesture
Restart=always
RestartSec=5

# Required for camera access
DeviceAllow=/dev/video*

[Install]
WantedBy=multi-user.target
